# deploy/alibaba/engram.service
[Unit]
# Engram hero (memory + sleep + viewer) as a systemd service on the ECS VM.
# Install: copy to /etc/systemd/system/, replace __USER__ and __DIR__, then
#   sudo systemctl enable --now engram
# (The chat agent / nanoclaw host is run separately — see docs/agent-and-deploy.md.)
Description=Engram (memory + viewer)
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
# nvm-installed node needs to be on PATH; adjust if node is elsewhere.
Environment=PATH=/home/__USER__/.nvm/versions/node/v20.19.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# `up` boots infra + viewer + sleep daemons in the background and returns; we keep
# the unit alive by tailing, and stop cleanly on shutdown.
ExecStart=/bin/bash -lc './engram.sh up && tail -F .data/run/viewer.log'
ExecStop=/bin/bash -lc './engram.sh down'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
